# systemd/pencil_test@.service
[Unit]
Description=Pencil-code test (master branch)

[Service]
Type=oneshot
WorkingDirectory=%I
EnvironmentFile=%I/pencil_test.conf

ExecStartPre=mkdir -p ../${output_dir_prev_name}
#NOTE: below is needed since otherwise the stdout file is not truncated; need to figure out a way to do that (simply file:->truncate: will not work, as that seems to truncate for each subcommand. Note that if I don't truncate, index_last.html will not be correctly generated by pencil-test. Ugh.
ExecStartPre=bash -c 'if test -e index.html; then mv index.html ../${output_dir_prev_name}/; fi'

ExecStart=bash -c '${pencil_home}/bin/pencil-test \
	--pencil-home=${pencil_home} \
	-N 15 --update --html --clean --local --use-pc_auto-test --coverage \
	--auto-test-options="--max-level=3 --script-tests=python --time-limit=5m \
		--auto-clean" \
	--log-dir=. --previous-dir=../${output_dir_prev_name} \
	> index.html'

